# .config/systemd/user/bing-get.timer
[Unit]
Description=Auto download bing wallpaper timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
