# .config/systemd/user/swww-next.timer
[Unit]
Description=Auto load next img by swww timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
